# client/etc/systemd/system/rs-backup-run.service
[Unit]
Description=Start backup via rs-backup-suite

[Service]
Type=simple
ExecStart=/usr/bin/rs-backup-run
