# server/etc/cron.weekly/rs-backup-rotate
#!/bin/sh
/usr/local/sbin/rs-rotate-cron weekly
