# server/etc/crontab
@daily    /usr/sbin/rs-rotate-cron daily
@weekly   /usr/sbin/rs-rotate-cron weekly
@monthly  /usr/sbin/rs-rotate-cron monthly
