# server/etc/cron.monthly/rs-backup-rotate
#!/bin/sh
/usr/sbin/rs-rotate-cron monthly
